# packages/nosaastrust/nosaastrust-1.0rc2.dev0.tar.gz/nosaastrust-1.0rc2.dev0/resources/nost-user.service.in
# SPDX-License-Identifier: BSD-2-Clause-Patent
# SPDX-FileCopyrightText: 2024 Syslinbit <contact@syslinbit.com>

[Unit]
Description=NoSaasTrust service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot

# If python venv is used to build binary :
ExecStart=@BINPATH@/notst @CONFIGPATH@/config.yml

[Install]
WantedBy=nost.timer
